# data/packagekit.service.in
[Unit]
Description=PackageKit Daemon
# PK doesn't know how to do anything on ostree-managed systems;
# currently the design is to have dedicated daemons like
# eos-updater and rpm-ostree, and gnome-software talks to those.
ConditionPathExists=!/run/ostree-booted
Wants=network-online.target

[Service]
Type=dbus
BusName=org.freedesktop.PackageKit
User=@PACKAGEKIT_USER@
# To enable verbose logging define following in sysconfig file:
# PACKAGEKIT_ARGS="--keep-environment --verbose"
# ZYPP_LOGFILE=-
EnvironmentFile=-/etc/sysconfig/packagekit
ExecStart=@libexecdir@/packagekitd $PACKAGEKIT_ARGS
